# pi/cortex_web.service
[Unit]
Description=CORTEX Web Service (FastAPI REST + WebSocket)
After=cortex_hub.service
BindsTo=cortex_hub.service

[Service]
Type=simple
User=pi
# Change this to the actual path where you cloned the project
WorkingDirectory=/home/pi/Cortex/pi
# This command assumes you have a virtualenv in the project root
ExecStart=/home/pi/Cortex/venv/bin/python cortex_web.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
